# bin/com.cogyard.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- EXAMPLE LaunchAgent — optional macOS "always-on" recipe (see docs/DEPLOY.md).
         Replace /ABSOLUTE/PATH/TO/cogyard and /Users/YOU before installing. -->
    <key>Label</key>
    <string>com.cogyard.serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/cogyard/bin/serve</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOU/.cogyard/logs/cogyard.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOU/.cogyard/logs/cogyard.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOU</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/local/bin</string>
    </dict>
</dict>
</plist>
